# aconfmgr/files/etc/profile.d/aurman.sh
#!/bin/sh

# Build packages in tmpfs
export AURDEST=/tmp/aurman

# Save packages in pacman cache
export PKGDEST=/var/cache/pacman/pkg
